# files/etc/cron.d/squidguard
07 03 * * * root /libexec/update-blacklists
