# files/etc/cron.d/squidguard
07 03 * * * proxy /libexec/update-blacklists
