# files/etc/cron.d/squidguard
SHELL=/bin/sh

07 03 * * * proxy /libexec/update-blacklists
